# services/mariadb/etc/supervisor/conf.d/mariadb.conf
[program:mariadb]
command=/setup/setup.sh && mysqld_safe
