# services/mariadb/etc/supervisor/conf.d/mariadb.conf
[program:mariadb]
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /bin/bash -c "/setup/setup.sh && mysqld_safe"
redirect_stderr=true
