# services/mariadb/etc/supervisor/conf.d/mariadb.conf
[program:mariadb]
command=/setup/setup.sh && /usr/bin/pidproxy /var/run/mysqld/mysqld.pid mysqld_safe
redirect_stderr=true
