# applications/kibana4/etc/supervisor/conf.d/kibana4.conf
[supervisord]
nodaemon=true

[program:kibana4]
command=/opt/kibana/bin/kibana
